# src/systemd/gnome-schema-sources.path
[Unit]
Description=Watch the GNOME schema sources

[Path]
PathChanged=%h/.local/share/gnome-shell/extensions
PathChanged=/usr/share/glib-2.0/schemas
Unit=update-compiled-schemas.service

[Install]
WantedBy=paths.target
